# templates/initd_script.erb
#!/bin/sh

# Autogenerated by foreman

### BEGIN INIT INFO
# Provides:             <%= name %>
# Required-Start:       $syslog
# Required-Stop:        $syslog
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    <%= name %>
# Description:          <%= description %>
### END INIT INFO

NAME=<%= Shellwords.escape name %>
DESC=<%= Shellwords.escape description %>
DAEMON=<%= Shellwords.escape script %>
DAEMON_USER=<%= Shellwords.escape user %>
PIDDIR=<%= Shellwords.escape pidfile.dirname.to_s %>
PIDFILE=<%= Shellwords.escape pidfile.to_s %>

test -x ${DAEMON} || exit 0
set -e
. /lib/lsb/init-functions


case "${1}" in
	start)
		log_daemon_msg "Starting ${DESC}" "${NAME}"
		mkdir -p ${PIDDIR}
		if (start-stop-daemon --start --make-pidfile --background --pidfile $PIDFILE --chuid $DAEMON_USER --startas $DAEMON -- <%= Shellwords.join arguments %>); then
			log_end_msg 0
		else
			log_end_msg 1
		fi
	;;
	stop)
		log_daemon_msg "Stopping ${DESC}" "${NAME}"
		if (start-stop-daemon --stop --oknodo --retry <%= Shellwords.escape timeout_schedule %> --pidfile $PIDFILE --chuid $DAEMON_USER); then
			log_end_msg 0
		else
			log_end_msg 1
		fi
	;;
	status)
		status_of_proc -p ${PIDFILE} ${DAEMON} ${NAME}
	;;
	restart)
		${0} stop
		${0} start
	;;
	*)
		echo "Usage: ${0} {start|stop|status|restart}" >&2
	;;
esac

exit 0
